# tests/system/config/supervisor/named.conf
[program:named]
# The -g flag runs the daemon in the foreground and forces logging to stderr.
# It allows to read the BIND 9 logs with docker logs.
directory=/etc/bind
command=/usr/sbin/named -t /chroot -u bind -g -c /etc/bind/named.conf
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
