# tests/system/config/supervisor/stork-server.conf
[program:stork-server]
command=rake run:server_debug HEADLESS=true UI_MODE=none
directory=/app
autostart = true
autorestart = unexpected
# Status 0 - expected program exit, 130 - SIGINT, CTRL+C
# Unfortunately, when the program is running through delve, as in this case,
# the actual exit code is suppressed and the delve exit code is returned.
exitcodes = 0,130
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
