# tests/system/config/supervisor/stork-server.conf
[program:stork-server]
command=rake run:server_debug HEADLESS=true UI_MODE=none
directory=/app
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
