# tests/system/config/supervisor/stork-server.conf
[program:stork-server]
command=/usr/bin/stork-server
autostart = true
autorestart = unexpected
# Status 0 - expected program exit, 130 - SIGINT, CTRL+C
exitcodes = 0,130
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
