# tests/system/config/supervisor/supervisor-profiler.conf
[program:supervisor-profiler]
command=python3 /usr/lib/supervisor-profiler.py /var/log/supervisor/performance-report
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
