# etc/isc-stork-server.initd
#!/sbin/openrc-run

cfgfile="/etc/stork/server.conf"
command="/usr/bin/stork-server"
command_args="--use-env-file"
command_background="yes"
command_user="stork-server"
output_log="/var/log/isc-stork-server.log"
pidfile="/run/isc-stork-server.pid"

depend() {
    use net
}

start_pre() {
    touch "${output_log}"
    chown "${command_user}" "${output_log}"
}
