# etc/stork.initd.in
#!/sbin/openrc-run

daemon=$(echo "${RC_SVCNAME##*-}")

cfgfile="/etc/stork/${daemon}.conf"
command="/usr/bin/stork-${daemon}"
command_args="--use-env-file"
command_background="yes"
command_user="stork-${daemon}"
output_log="/var/log/${RC_SVCNAME}.log"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    use net
}

start_pre() {
    touch "${output_log}"
    chown "${command_user}" "${output_log}"
}
