# etc/hooks/rpm/isc-stork-server.prerm
#!/bin/sh

if [ "$(systemctl is-system-running)" -eq 0 ]; then
    systemctl disable isc-stork-server
    systemctl stop isc-stork-server
fi
