# docker/config/supervisor/named.conf
[program:named]
# The -g flag runs the daemon in the foreground and forces logging to stderr.
# It allows to read the BIND 9 logs with docker logs. BIND needs writeable
# working directory. Real deployments would use /var/cache/bind or similar,
# but for demo docker image, /etc/bind is good enough.
directory=/etc/bind
command=/usr/sbin/named -g -c /etc/bind/named.conf
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
