# docker/config/supervisor/supervisor_monitor.conf
[program:supervisor_monitor]
command=python3 /usr/lib/supervisor_monitor.py /var/log/supervisor/performance-report 2000
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
